# docker/dev/supervisord/sshd.conf
[program:sshd]
directory=/home/user
command=/usr/sbin/sshd -D